# minecraftd.service.in
[Unit]
Description=@GAME@ Server
After=local-fs.target network.target multi-user.target

[Service]
Type=forking
ExecStart=/usr/local/bin/@INAME@ start
ExecStop=/usr/local/bin/@INAME@ stop
User=@GAME_USER@
Group=@GAME_USER@

[Install]
WantedBy=multi-user.target
